# ansible/resources/playlist-manager-frontend.service
[Unit]
Description=Playlist Manager Frontend

[Service]
User=ec2-user
WorkingDirectory=/opt/playlist-manager/frontend
ExecStart=/usr/bin/npm run dev